# system_scripts/install_service/web_server.service
[Unit]
Description=Node Web Server
After=network-online.target
Wants=network-online.target

[Service]
Restart=always
ExecStart=

[Install]
WantedBy=default.target
